# mveg.cron
*/1 * * * * root sh /usr/local/mveg/mveg.sh > /dev/null 2>&1
