# mveg.cron
*/1 * * * * root . /etc/bashrc; sh /usr/local/mveg/mveg.sh > /dev/null 2>&1
